# monitor/systemd/vps-monitor-peers.timer
[Unit]
Description=Check WireGuard peers every 30 seconds

[Timer]
OnBootSec=60
OnUnitActiveSec=30s

[Install]
WantedBy=timers.target
